# sickrage.ini
[program:sickrage]
command=/SickRage-master/SickBeard.py --datadir=/config
