# sickrage.ini
[program:sickrage]
command=/SickRage/SickBeard.py --datadir=/config
