# sickrage.ini
[program:sickrage]
user=
command=/SickRage/SickBeard.py --datadir=/config
